# bin/scripts and services to launch applications/services/app-services.service
[Unit]
Description = launch services application with external configuration files.

[Service]
Type=simple
ExecStart=/bin/bash /path/to/script/app-services.sh
TimeoutStartSec=0

[Install]
WantedBy = multi-user.target
